# deployment/celery.service
[Unit]
Description=Blogify Celery Worker
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/blogify
ExecStart=/var/www/blogify/venv/bin/celery -A blogify worker --loglevel=info
Restart=always
RestartSec=5
Environment="PATH=/var/www/blogify/venv/bin"

[Install]
WantedBy=multi-user.target 