# deployment/crontab
# Auto-backup database daily at 3:00 AM
0 3 * * * /var/www/blogify/venv/bin/python /var/www/blogify/manage.py backup_database

# Auto-update system daily at 4:00 AM (will only pull updates if new commits are available)
0 4 * * * /var/www/blogify/deployment/update.sh > /var/log/blogify-update.log 2>&1

# Check and restart services if they're down, runs every 5 minutes
*/5 * * * * systemctl is-active --quiet gunicorn || systemctl start gunicorn
*/5 * * * * systemctl is-active --quiet celery || systemctl start celery
*/5 * * * * systemctl is-active --quiet celerybeat || systemctl start celerybeat
*/5 * * * * systemctl is-active --quiet redis-server || systemctl start redis-server 